# scripts/com.dotfiles.autocommit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles.autocommit</string>

    <key>StartInterval</key>
    <integer>604800</integer> <!-- Run every 7 days -->

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>~/.dotfiles/autocommit.sh</string>
    </array>

    <key>StandardOutPath</key>
    <string>~/Library/Logs/com.dotfiles.autocommit.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/com.dotfiles.autocommit.err</string>
</dict>
</plist>
